# distreboot/packaging/fedora/stasher-distreboot.service.in
# Copyright 2012 Double Precision, Inc.
# See COPYING for distribution information.

[Unit]
Description=A distributed reboot service
Wants=stasher.service
After=stasher.service

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
ExecStart=@sbindir@/distreboot --start
ExecStop=@sbindir@/distreboot --stop
KillMode=process

