# distreboot/packaging/freebsd9/databases/stasher-distreboot/files/distreboot.in
#! /bin/sh
#
#

# REQUIRE: stasherd
# PROVIDE: distreboot
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable distreboot:
#
#distreboot_enable="YES"

. /etc/rc.subr

name="distreboot"
rcvar=`set_rcvar`

load_rc_config $name

start_cmd="distreboot_start"
stop_cmd="distreboot_stop"

prefix="%%PREFIX%%"

distreboot_start() {

    $prefix/sbin/distreboot --start --daemon `cat $prefix/etc/distrebootrc 2>/dev/null`
}

distreboot_stop() {

    $prefix/sbin/distreboot --stop
}
run_rc_command "$1"
